# com.noroshi.github-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.noroshi.github-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tat/Library/Mobile Documents/iCloud~md~obsidian/Documents/NOROSHI/auto-sync-github.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>30</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/tat/Library/Mobile Documents/iCloud~md~obsidian/Documents/NOROSHI/github-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tat/Library/Mobile Documents/iCloud~md~obsidian/Documents/NOROSHI/github-sync-error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/tat/Library/Mobile Documents/iCloud~md~obsidian/Documents/NOROSHI</string>
</dict>
</plist> 